# go/deploy/spire/contrib/systemd/spire-agent.service.d/auto-join.conf
[Service]
# This file provides the join token for automatic agent registration
Environment="UNKEY_SPIRE_JOIN_TOKEN="
Environment="UNKEY_SPIRE_TRUST_DOMAIN=development.unkey.cloud"
