# go/deploy/spire/contrib/systemd/spire-agent.service
[Unit]
Description=SPIRE Agent
Documentation=https://spiffe.io/docs/latest/
After=network.target
Wants=network-online.target
# Server connection will retry if not available

[Service]
Type=simple
User=root
Group=root

# StateDirectory creates persistent directories under /var/lib
StateDirectory=spire/agent spire/agent/data spire/agent/keys
StateDirectoryMode=0755

ExecStartPre=/bin/bash -c 'chmod 700 /var/lib/spire/agent/keys'
ExecStart=/opt/spire/bin/spire-agent-wrapper.sh

ExecStop=/bin/kill -SIGTERM $MAINPID
TimeoutStopSec=30

Restart=on-failure
RestartSec=5
TimeoutStartSec=30s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=spire-agent

NoNewPrivileges=true
PrivateTmp=no
ProtectSystem=no
ProtectHome=no
ReadWritePaths=/var/lib/spire/agent
ReadOnlyPaths=/usr/bin /usr/local/bin /opt
# Keep some basic protections
RestrictRealtime=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
SystemCallFilter=@system-service

# SupplementaryGroups would be added here if services run under separate users

Environment="UNKEY_SPIRE_TRUST_DOMAIN=development.unkey.cloud"
Environment="UNKEY_SPIRE_LOG_LEVEL=INFO"
Environment="UNKEY_SPIRE_SERVER_URL=https://localhost:8085"
# Set this in environment-specific drop-ins for security
# Environment="UNKEY_SPIRE_JOIN_TOKEN=your-long-lived-token-here"
# e.g., /etc/systemd/system/spire-agent.service.d/environment.conf

[Install]
WantedBy=multi-user.target
