# go/deploy/spire/contrib/systemd/spire-server-fixed.service
[Unit]
Description=SPIRE Server
Documentation=https://spiffe.io/docs/latest/
After=network.target
Wants=network-online.target

[Service]
Type=simple
# AIDEV-NOTE: Run as root for full system access
User=root
Group=root

# AIDEV-NOTE: Systemd-managed directory creation
# RuntimeDirectory creates /run/spire (cleared on reboot)
RuntimeDirectory=spire spire/sockets
RuntimeDirectoryMode=0755
RuntimeDirectoryPreserve=yes  # Prevent cleanup while service is running

# StateDirectory creates /var/lib/spire (persistent)
StateDirectory=spire/server spire/server/data spire/server/keys
StateDirectoryMode=0755

# AIDEV-NOTE: Additional permission setup for keys directory
ExecStartPre=/bin/bash -c 'chmod 700 /var/lib/spire/server/keys || true'

# AIDEV-NOTE: Main server process with config path
ExecStart=/opt/spire/bin/spire-server run -config /etc/spire/server/server.conf

# AIDEV-NOTE: Graceful shutdown
ExecStop=/bin/kill -SIGTERM $MAINPID
TimeoutStopSec=30

Restart=on-failure
RestartSec=5
TimeoutStartSec=30s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=spire-server

# AIDEV-SECURITY: Adjusted security settings for SPIRE's needs
NoNewPrivileges=true
PrivateTmp=true
# Changed from 'full' to 'strict' with explicit ReadWritePaths
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/spire/server /run/spire
# Allow reading config
ReadOnlyPaths=/etc/spire
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
# Removed MemoryDenyWriteExecute as it might interfere with Go runtime
# MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
SystemCallFilter=@system-service

# AIDEV-NOTE: Environment configuration
# Trust domain must be set per environment
Environment="UNKEY_SPIRE_TRUST_DOMAIN=development.unkey.app"
Environment="UNKEY_SPIRE_LOG_LEVEL=INFO"
Environment="UNKEY_SPIRE_DB_TYPE=sqlite3"
Environment="UNKEY_SPIRE_DB_CONNECTION=/var/lib/spire/server/data/datastore.sqlite3"

[Install]
WantedBy=multi-user.target