# go/deploy/spire/contrib/systemd/spire-server.service
[Unit]
Description=SPIRE Server
Documentation=https://spiffe.io/docs/latest/
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root

# RuntimeDirectory creates /run/spire (cleared on reboot) - keeping for backward compatibility
RuntimeDirectory=spire
RuntimeDirectoryMode=0755
RuntimeDirectoryPreserve=yes
# StateDirectory creates /var/lib/spire (persistent)
StateDirectory=spire/server spire/server/data spire/server/keys
StateDirectoryMode=0755
# ConfigurationDirectory would create /etc/spire but we manage this separately
# ConfigurationDirectory=spire/server

ExecStartPre=/bin/bash -c 'chmod 700 /var/lib/spire/server/keys || true'

ExecStart=/opt/spire/bin/spire-server run -config /etc/spire/server/server.conf

ExecStop=/bin/kill -SIGTERM $MAINPID
TimeoutStopSec=30

Restart=on-failure
RestartSec=5
TimeoutStartSec=30s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=spire-server

# Most security features disabled to prevent namespace isolation
NoNewPrivileges=true
PrivateTmp=no
ProtectSystem=no
ProtectHome=no
ReadWritePaths=/var/lib/spire/server /run/spire
# Keep some basic protections
RestrictRealtime=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# Trust domain must be set per environment
Environment="UNKEY_SPIRE_TRUST_DOMAIN=development.unkey.cloud"
Environment="UNKEY_SPIRE_LOG_LEVEL=INFO"
Environment="UNKEY_SPIRE_DB_TYPE=sqlite3"
Environment="UNKEY_SPIRE_DB_CONNECTION=/var/lib/spire/server/data/datastore.sqlite3"
# e.g., /etc/systemd/system/spire-server.service.d/environment.conf
# Optional: Set UNKEY_SPIRE_TRUST_BUNDLE if using external trust bundle

[Install]
WantedBy=multi-user.target
