# go/deploy/docker/systemd/assetmanagerd.service
[Unit]
Description=AssetManagerd VM Asset Management Service (Development)
Documentation=https://github.com/unkeyed/unkey/go/assetmanagerd
After=network.target
Wants=network.target

[Service]
Type=simple
# Run as root for development simplicity
User=root
Group=root
# Load development environment variables
EnvironmentFile=/etc/default/unkey-deploy
# Service will start in the required directories
ExecStart=/usr/local/bin/assetmanagerd
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=assetmanagerd

# Development-specific environment overrides
Environment=UNKEY_ASSETMANAGERD_TLS_MODE=disabled
Environment=UNKEY_ASSETMANAGERD_OTEL_ENABLED=false
Environment=UNKEY_ASSETMANAGERD_PORT=8083
Environment=UNKEY_ASSETMANAGERD_ADDRESS=0.0.0.0

# Storage configuration for development
Environment=UNKEY_ASSETMANAGERD_STORAGE_BACKEND=local
Environment=UNKEY_ASSETMANAGERD_LOCAL_STORAGE_PATH=/opt/builderd/rootfs
Environment=UNKEY_ASSETMANAGERD_DATABASE_PATH=/opt/assetmanagerd/assets.db
Environment=UNKEY_ASSETMANAGERD_CACHE_DIR=/opt/assetmanagerd/cache

# Garbage collection
Environment=UNKEY_ASSETMANAGERD_GC_ENABLED=true
Environment=UNKEY_ASSETMANAGERD_GC_INTERVAL=1h
Environment=UNKEY_ASSETMANAGERD_GC_MAX_AGE=24h

# Builderd integration configuration
Environment=UNKEY_ASSETMANAGERD_BUILDERD_ENABLED=true
Environment=UNKEY_ASSETMANAGERD_BUILDERD_ENDPOINT=http://localhost:8082
Environment=UNKEY_ASSETMANAGERD_BUILDERD_TIMEOUT=30m
Environment=UNKEY_ASSETMANAGERD_BUILDERD_AUTO_REGISTER=true
Environment=UNKEY_ASSETMANAGERD_BUILDERD_MAX_RETRIES=3
Environment=UNKEY_ASSETMANAGERD_BUILDERD_RETRY_DELAY=5s

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Basic security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
