# go/deploy/docker/systemd/billaged.service
[Unit]
Description=Billaged VM Usage Billing Service (Development)
Documentation=https://github.com/unkeyed/unkey/go/deploy/billaged
After=network.target
Wants=network.target

[Service]
Type=simple
User=billaged
Group=billaged
# Load development environment variables
EnvironmentFile=/etc/default/unkey-deploy
# Service will start in the required directories
ExecStart=/usr/local/bin/billaged
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=billaged

# Development-specific environment overrides
Environment=UNKEY_BILLAGED_TLS_MODE=disabled
Environment=UNKEY_BILLAGED_OTEL_ENABLED=false
Environment=UNKEY_BILLAGED_PORT=8081
Environment=UNKEY_BILLAGED_ADDRESS=0.0.0.0
Environment=UNKEY_BILLAGED_AGGREGATION_INTERVAL=60s

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target