# go/deploy/docker/systemd/builderd.service
[Unit]
Description=Builderd Multi-Tenant Build Service (Development)
Documentation=https://github.com/unkeyed/unkey/go/deploy/builderd
After=network.target
Wants=network.target

[Service]
Type=simple
# Run as root for Docker access
User=root
Group=root
# Load development environment variables
EnvironmentFile=/etc/default/unkey-deploy
# Service will start in the required directories
ExecStart=/usr/local/bin/builderd
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=builderd

# Development-specific environment overrides
Environment=UNKEY_BUILDERD_TLS_MODE=disabled
Environment=UNKEY_BUILDERD_OTEL_ENABLED=false
Environment=UNKEY_BUILDERD_PORT=8082
Environment=UNKEY_BUILDERD_ADDRESS=0.0.0.0

# Build configuration for development
Environment=UNKEY_BUILDERD_MAX_CONCURRENT_BUILDS=3
Environment=UNKEY_BUILDERD_BUILD_TIMEOUT=15m
Environment=UNKEY_BUILDERD_SCRATCH_DIR=/opt/builderd/scratch
Environment=UNKEY_BUILDERD_ROOTFS_OUTPUT_DIR=/opt/builderd/rootfs
Environment=UNKEY_BUILDERD_WORKSPACE_DIR=/opt/builderd/workspace

# Storage configuration
Environment=UNKEY_BUILDERD_STORAGE_BACKEND=local
Environment=UNKEY_BUILDERD_STORAGE_RETENTION_DAYS=7

# Database configuration
Environment=UNKEY_BUILDERD_DATABASE_TYPE=sqlite
Environment=UNKEY_BUILDERD_DATABASE_DATA_DIR=/opt/builderd/data

# Docker configuration
Environment=UNKEY_BUILDERD_DOCKER_MAX_IMAGE_SIZE_GB=2

# Tenant isolation (disabled for development)
Environment=UNKEY_BUILDERD_TENANT_ISOLATION_ENABLED=false

# AssetManagerd integration
Environment=UNKEY_BUILDERD_ASSETMANAGER_ENABLED=true
Environment=UNKEY_BUILDERD_ASSETMANAGER_ENDPOINT=http://localhost:8083

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Shutdown configuration
TimeoutStopSec=30
KillMode=mixed

[Install]
WantedBy=multi-user.target