# go/deploy/docker/systemd/metald.service
[Unit]
Description=Metald VM Management Service (Development)
Documentation=https://github.com/unkeyed/unkey/go/deploy/metald
After=network.target assetmanagerd.service billaged.service
Wants=network.target
# Start after assetmanagerd and billaged for proper service order

[Service]
Type=simple
# Run as root for Docker access
User=root
Group=root
# Load development environment variables
EnvironmentFile=/etc/default/unkey-deploy
# Service will start in the required directories
ExecStart=/usr/local/bin/metald
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=metald

# Development-specific environment overrides
Environment=UNKEY_METALD_TLS_MODE=disabled
Environment=UNKEY_METALD_OTEL_ENABLED=false
Environment=UNKEY_METALD_PORT=8080
Environment=UNKEY_METALD_ADDRESS=0.0.0.0

# Docker backend configuration
Environment=UNKEY_METALD_BACKEND=docker
Environment=UNKEY_METALD_DOCKER_HOST=unix:///var/run/docker.sock

# Process Manager Configuration
Environment=UNKEY_METALD_SOCKET_DIR=/opt/metald/sockets
Environment=UNKEY_METALD_LOG_DIR=/opt/metald/logs
Environment=UNKEY_METALD_MAX_PROCESSES=100

# Billing Configuration
Environment=UNKEY_METALD_BILLING_ENABLED=true
Environment=UNKEY_METALD_BILLING_ENDPOINT=http://localhost:8081
Environment=UNKEY_METALD_BILLING_MOCK_MODE=false

# AssetManager Configuration
Environment=UNKEY_METALD_ASSETMANAGER_ENABLED=true
Environment=UNKEY_METALD_ASSETMANAGER_ENDPOINT=http://localhost:8083
Environment=UNKEY_METALD_ASSETMANAGER_CACHE_DIR=/opt/metald/assets

# Network configuration for Docker backend
Environment=UNKEY_METALD_DOCKER_NETWORK=bridge
Environment=UNKEY_METALD_DOCKER_PORT_RANGE_MIN=30000
Environment=UNKEY_METALD_DOCKER_PORT_RANGE_MAX=40000

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target