# go/apps/billaged/contrib/systemd/billaged.service
[Unit]
Description=Billaged VM Usage Billing Service
Documentation=https://github.com/unkeyed/unkey/go/apps/billaged
After=network.target spire-agent.service assetmanagerd.service metald.service
Wants=network.target
Requires=spire-agent.service assetmanagerd.service metald.service

[Service]
Type=simple
User=billaged
Group=billaged
# AIDEV-NOTE: WorkingDirectory removed - not needed for billaged
# Create required directories (+ prefix runs as root)
ExecStartPre=+/usr/bin/mkdir -p /opt/billaged
ExecStartPre=+/usr/bin/mkdir -p /var/log/billaged
# Set ownership for service directories
ExecStartPre=+/usr/bin/chown -R billaged:billaged /opt/billaged
ExecStartPre=+/usr/bin/chown -R billaged:billaged /var/log/billaged
ExecStart=/usr/local/bin/billaged
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=billaged

# Environment variables
Environment=UNKEY_BILLAGED_PORT=8081
Environment=UNKEY_BILLAGED_ADDRESS=0.0.0.0
Environment=UNKEY_BILLAGED_AGGREGATION_INTERVAL=60s

# OpenTelemetry Configuration (enabled for production)
Environment=UNKEY_BILLAGED_OTEL_ENABLED=true
Environment=UNKEY_BILLAGED_OTEL_SERVICE_NAME=billaged
Environment=UNKEY_BILLAGED_OTEL_SERVICE_VERSION=0.0.1
Environment=UNKEY_BILLAGED_OTEL_SAMPLING_RATE=1.0
Environment=UNKEY_BILLAGED_OTEL_ENDPOINT=localhost:4318
Environment=UNKEY_BILLAGED_OTEL_PROMETHEUS_ENABLED=true
Environment=UNKEY_BILLAGED_OTEL_PROMETHEUS_PORT=9465

# TLS/SPIFFE configuration (REQUIRED)
# AIDEV-BUSINESS_RULE: mTLS is required for secure inter-service communication
Environment=UNKEY_BILLAGED_TLS_MODE=spiffe
Environment=UNKEY_BILLAGED_SPIFFE_SOCKET=/var/lib/spire/agent/agent.sock

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
