# go/deploy/assetmanagerd/contrib/systemd/assetmanagerd.service
[Unit]
Description=AssetManagerd VM Asset Management Service
Documentation=https://github.com/unkeyed/unkey/go/deploy/assetmanagerd
After=network.target spire-agent.service
Wants=network.target
Requires=spire-agent.service

[Service]
Type=simple
# Running as root for cross-service file access
User=root
Group=root
# AIDEV-NOTE: WorkingDirectory removed - not needed for assetmanagerd
# Create required directories (+ prefix runs as root before dropping privileges)
ExecStartPre=+/usr/bin/mkdir -p /opt/assetmanagerd/cache
ExecStartPre=+/usr/bin/mkdir -p /opt/vm-assets
ExecStartPre=+/usr/bin/mkdir -p /var/log/assetmanagerd
# No ownership changes needed when running as root
ExecStart=/usr/local/bin/assetmanagerd
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=assetmanagerd

# Environment variables
Environment=UNKEY_ASSETMANAGERD_PORT=8083
Environment=UNKEY_ASSETMANAGERD_ADDRESS=0.0.0.0

# Storage configuration
Environment=UNKEY_ASSETMANAGERD_STORAGE_BACKEND=local
Environment=UNKEY_ASSETMANAGERD_LOCAL_STORAGE_PATH=/opt/builderd/rootfs
Environment=UNKEY_ASSETMANAGERD_DATABASE_PATH=/opt/assetmanagerd/assets.db
Environment=UNKEY_ASSETMANAGERD_CACHE_DIR=/opt/assetmanagerd/cache

# Garbage collection
Environment=UNKEY_ASSETMANAGERD_GC_ENABLED=true
Environment=UNKEY_ASSETMANAGERD_GC_INTERVAL=1h
Environment=UNKEY_ASSETMANAGERD_GC_MAX_AGE=168h

# OpenTelemetry Configuration
Environment=UNKEY_ASSETMANAGERD_OTEL_ENABLED=true
Environment=UNKEY_ASSETMANAGERD_OTEL_SERVICE_NAME=assetmanagerd
Environment=UNKEY_ASSETMANAGERD_OTEL_SERVICE_VERSION=0.2.0
Environment=UNKEY_ASSETMANAGERD_OTEL_SAMPLING_RATE=1.0
Environment=UNKEY_ASSETMANAGERD_OTEL_ENDPOINT=localhost:4318
Environment=UNKEY_ASSETMANAGERD_OTEL_PROMETHEUS_ENABLED=true
Environment=UNKEY_ASSETMANAGERD_OTEL_PROMETHEUS_PORT=9467

# TLS/SPIFFE configuration (REQUIRED)
# AIDEV-BUSINESS_RULE: mTLS is required for secure inter-service communication
Environment=UNKEY_ASSETMANAGERD_TLS_MODE=spiffe
Environment=UNKEY_ASSETMANAGERD_SPIFFE_SOCKET=/var/lib/spire/agent/agent.sock

# Builderd integration configuration
# AIDEV-NOTE: Enables automatic rootfs building when assets don't exist
Environment=UNKEY_ASSETMANAGERD_BUILDERD_ENABLED=true
Environment=UNKEY_ASSETMANAGERD_BUILDERD_ENDPOINT=https://localhost:8082
Environment=UNKEY_ASSETMANAGERD_BUILDERD_TIMEOUT=30m
Environment=UNKEY_ASSETMANAGERD_BUILDERD_AUTO_REGISTER=true
Environment=UNKEY_ASSETMANAGERD_BUILDERD_MAX_RETRIES=3
Environment=UNKEY_ASSETMANAGERD_BUILDERD_RETRY_DELAY=5s

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Basic security settings (removed strict namespace protection)
# AIDEV-NOTE: Namespace protection removed to simplify deployment
# The service runs as root for filesystem operations
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
