# go/deploy/assetmanagerd/contrib/systemd/assetmanagerd.service
[Unit]
Description=AssetManagerd VM Asset Management Service
Documentation=https://github.com/unkeyed/unkey/go/deploy/assetmanagerd
After=network.target
Wants=network.target

[Service]
Type=simple
User=assetmanagerd
Group=assetmanagerd
WorkingDirectory=/opt/assetmanagerd
ExecStart=/usr/local/bin/assetmanagerd
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=assetmanagerd

# Environment variables
Environment=UNKEY_ASSETMANAGERD_PORT=8083
Environment=UNKEY_ASSETMANAGERD_ADDRESS=0.0.0.0

# Storage configuration
Environment=UNKEY_ASSETMANAGERD_STORAGE_BACKEND=local
Environment=UNKEY_ASSETMANAGERD_LOCAL_STORAGE_PATH=/opt/vm-assets
Environment=UNKEY_ASSETMANAGERD_DATABASE_PATH=/opt/assetmanagerd/assets.db
Environment=UNKEY_ASSETMANAGERD_CACHE_DIR=/opt/assetmanagerd/cache

# Garbage collection
Environment=UNKEY_ASSETMANAGERD_GC_ENABLED=true
Environment=UNKEY_ASSETMANAGERD_GC_INTERVAL=1h
Environment=UNKEY_ASSETMANAGERD_GC_MAX_AGE=168h

# OpenTelemetry Configuration
Environment=UNKEY_ASSETMANAGERD_OTEL_ENABLED=true
Environment=UNKEY_ASSETMANAGERD_OTEL_SERVICE_NAME=assetmanagerd
Environment=UNKEY_ASSETMANAGERD_OTEL_SERVICE_VERSION=0.0.1
Environment=UNKEY_ASSETMANAGERD_OTEL_SAMPLING_RATE=1.0
Environment=UNKEY_ASSETMANAGERD_OTEL_ENDPOINT=localhost:4318
Environment=UNKEY_ASSETMANAGERD_OTEL_PROMETHEUS_ENABLED=true
Environment=UNKEY_ASSETMANAGERD_OTEL_PROMETHEUS_PORT=9466

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security settings
NoNewPrivileges=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/assetmanagerd /opt/vm-assets

[Install]
WantedBy=multi-user.target