# go/deploy/builderd/contrib/systemd/builderd.service
[Unit]
Description=Builderd Multi-Tenant Build Service
Documentation=https://github.com/unkeyed/unkey/go/deploy/builderd
After=network.target
Wants=network.target

[Service]
Type=simple
User=builderd
Group=builderd
WorkingDirectory=/opt/builderd
ExecStart=/usr/local/bin/builderd
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=builderd

# Core service configuration
Environment=UNKEY_BUILDERD_PORT=8082
Environment=UNKEY_BUILDERD_ADDRESS=0.0.0.0

# Build configuration
Environment=UNKEY_BUILDERD_MAX_CONCURRENT_BUILDS=5
Environment=UNKEY_BUILDERD_BUILD_TIMEOUT=15m
Environment=UNKEY_BUILDERD_SCRATCH_DIR=/opt/builderd/scratch
Environment=UNKEY_BUILDERD_ROOTFS_OUTPUT_DIR=/opt/builderd/rootfs
Environment=UNKEY_BUILDERD_WORKSPACE_DIR=/opt/builderd/workspace

# Storage configuration
Environment=UNKEY_BUILDERD_STORAGE_BACKEND=local
Environment=UNKEY_BUILDERD_STORAGE_RETENTION_DAYS=30

# Database configuration
Environment=UNKEY_BUILDERD_DATABASE_TYPE=sqlite
Environment=UNKEY_BUILDERD_DATABASE_DATA_DIR=/opt/builderd/data

# Docker configuration
Environment=UNKEY_BUILDERD_DOCKER_MAX_IMAGE_SIZE_GB=5

# Tenant isolation
Environment=UNKEY_BUILDERD_TENANT_ISOLATION_ENABLED=true

# OpenTelemetry Configuration (enabled for production)
Environment=UNKEY_BUILDERD_OTEL_ENABLED=true
Environment=UNKEY_BUILDERD_OTEL_SERVICE_NAME=builderd
Environment=UNKEY_BUILDERD_OTEL_SERVICE_VERSION=0.0.1
Environment=UNKEY_BUILDERD_OTEL_SAMPLING_RATE=1.0
Environment=UNKEY_BUILDERD_OTEL_ENDPOINT=localhost:4318
Environment=UNKEY_BUILDERD_OTEL_PROMETHEUS_ENABLED=true
Environment=UNKEY_BUILDERD_OTEL_HIGH_CARDINALITY_ENABLED=false

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target