# go/deploy/builderd/contrib/systemd/builderd.service
[Unit]
Description=Builderd Multi-Tenant Build Service
Documentation=https://github.com/unkeyed/unkey/go/deploy/builderd
After=network.target spire-agent.service assetmanagerd.service
Wants=network.target
Requires=spire-agent.service assetmanagerd.service

[Service]
Type=simple
# Running as root for filesystem operations
User=root
Group=root
# AIDEV-NOTE: WorkingDirectory removed - not needed for builderd
# Create required directories (+ prefix runs as root)
ExecStartPre=+/usr/bin/mkdir -p /opt/builderd/scratch
ExecStartPre=+/usr/bin/mkdir -p /opt/builderd/rootfs
ExecStartPre=+/usr/bin/mkdir -p /opt/builderd/workspace
ExecStartPre=+/usr/bin/mkdir -p /opt/builderd/data
ExecStartPre=+/usr/bin/mkdir -p /var/log/builderd
# No ownership changes needed when running as root
ExecStart=/usr/local/bin/builderd
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=builderd

# Core service configuration
Environment=UNKEY_BUILDERD_PORT=8082
Environment=UNKEY_BUILDERD_ADDRESS=0.0.0.0

# Build configuration
Environment=UNKEY_BUILDERD_MAX_CONCURRENT_BUILDS=5
Environment=UNKEY_BUILDERD_BUILD_TIMEOUT=15m
Environment=UNKEY_BUILDERD_SCRATCH_DIR=/opt/builderd/scratch
Environment=UNKEY_BUILDERD_ROOTFS_OUTPUT_DIR=/opt/builderd/rootfs
Environment=UNKEY_BUILDERD_WORKSPACE_DIR=/opt/builderd/workspace
Environment=UNKEY_BUILDERD_USE_PIPELINE_EXECUTOR=false

# Storage configuration
Environment=UNKEY_BUILDERD_STORAGE_BACKEND=local
Environment=UNKEY_BUILDERD_STORAGE_RETENTION_DAYS=30

# Database configuration
Environment=UNKEY_BUILDERD_DATABASE_TYPE=sqlite
Environment=UNKEY_BUILDERD_DATABASE_DATA_DIR=/opt/builderd/data

# Docker configuration
Environment=UNKEY_BUILDERD_DOCKER_MAX_IMAGE_SIZE_GB=5

# Tenant isolation
Environment=UNKEY_BUILDERD_TENANT_ISOLATION_ENABLED=true

# AssetManagerd integration
Environment=UNKEY_BUILDERD_ASSETMANAGER_ENABLED=true
Environment=UNKEY_BUILDERD_ASSETMANAGER_ENDPOINT=https://localhost:8083

# OpenTelemetry Configuration (enabled for production)
Environment=UNKEY_BUILDERD_OTEL_ENABLED=true
Environment=UNKEY_BUILDERD_OTEL_SERVICE_NAME=builderd
Environment=UNKEY_BUILDERD_OTEL_SERVICE_VERSION=0.0.1
Environment=UNKEY_BUILDERD_OTEL_SAMPLING_RATE=1.0
Environment=UNKEY_BUILDERD_OTEL_ENDPOINT=localhost:4318
Environment=UNKEY_BUILDERD_OTEL_PROMETHEUS_ENABLED=true
Environment=UNKEY_BUILDERD_OTEL_PROMETHEUS_PORT=9466
Environment=UNKEY_BUILDERD_OTEL_PROMETHEUS_INTERFACE=127.0.0.1
Environment=UNKEY_BUILDERD_OTEL_HIGH_CARDINALITY_ENABLED=false

# TLS/SPIFFE configuration (REQUIRED)
# AIDEV-BUSINESS_RULE: mTLS is required for secure inter-service communication
Environment=UNKEY_BUILDERD_TLS_MODE=spiffe
Environment=UNKEY_BUILDERD_SPIFFE_SOCKET=/var/lib/spire/agent/agent.sock

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Shutdown configuration
# AIDEV-NOTE: Give builderd time to finish builds and shutdown gracefully
TimeoutStopSec=30
KillMode=mixed

[Install]
WantedBy=multi-user.target
